# stocks/DOW/scripts/crontabs
###### NEEDS EDITING ######
# m     h    dom mon dow   command
  0     5     *   *  1-5   `which ruby` #{root}/DOW/scripts/setup.rb

 29-59   6    *   *  1-5   `which ruby` #{root}/DOW/scripts/dow.rb -save > /tmp/stock.log 2>&1 &
  *     7-12  *   *  1-5   `which ruby` #{root}/DOW/scripts/dow.rb -save > /tmp/stock.log 2>&1 &
 0-10    13   *   *  1-5   `which ruby` #{root}/DOW/scripts/dow.rb -save > /tmp/stock.log 2>&1 &

  15    13    *   *  1-5   `which ruby` #{root}/DOW/scripts/analyze.rb -save
